# stt-citrinet/rootfs/etc/s6-overlay/s6-rc.d/stt-citrinet/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start stt-citrinet service
# ==============================================================================
cd /usr/src

flags=()

if bashio::config.true 'use_max_score'; then
    flags+=('--max-score' "$(bashio::config 'max_score')")
fi

if bashio::config.true 'debug_logging'; then
    flags+=('--debug')
fi

exec python3 app.py \
    --uri 'tcp://0.0.0.0:10300' \
    --model "$(bashio::config 'model')" \
    --http-host '0.0.0.0' \
    --http-port 5000 \
    --sentences /data/sentences.txt \
    --train-dir /data/train \
    --cache-dir /data/cache \
    ${flags[@]}
