# wake-phonmatchnet/rootfs/etc/s6-overlay/s6-rc.d/wake-phonmatchnet/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start wake-phonmatchnet service
# ==============================================================================
cd /usr/src

flags=()

if bashio::config.true 'debug_logging'; then
    flags+=('--debug')
fi

exec python3 app.py \
    --uri 'tcp://0.0.0.0:10400' \
    --wake-word "$(bashio::config 'wake_word')" \
    --threshold "$(bashio::config 'threshold')" \
    --trigger-count "$(bashio::config 'trigger_count')" \
    --checkpoint ./model \
    ${flags[@]}
