# canary-asr/rootfs/etc/s6-overlay/s6-rc.d/canary-asr/finish
#!/command/with-contenv bashio
# finish script